# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:btfs]
command=/usr/bin/btfs --api /ip4/0.0.0.0/tcp/5001 daemon
